# config/gunicorn.conf
description "yaari gunicorn"
start on runlevel [2345]
stop on runlevel [06]
#respawn
#respawn limit 10 5
exec su - ubuntu -c "bash /usr/local/bin/run_gunicorn_yaari.sh"
